# infrastructure/systemd/storeai-web.service.example
[Unit]
Description=StoreAI web (Next.js)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User={$SERVICE_USER}
WorkingDirectory={$REPO_DIR}
Environment=NODE_ENV=production
Environment=HOST=localhost
Environment=PORT={$STOREAI_PORT}
Environment=PATH={$SERVICE_PATH}
ExecStart={$PNPM_BIN} start
# Always restart — the app should survive crashes, OOMs, and DB/Redis blips.
Restart=always
RestartSec=5
# If it keeps crashing fast, back off. StartLimitBurst above caps the rate.
TimeoutStartSec=60
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
# Basic hardening (tight enough to catch path-traversal writes, loose enough
# for a Node app that reads node_modules + writes .next/cache)
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
