# infrastructure/systemd/storeai-web.service.example
[Unit]
Description=StoreAI web (Next.js)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User={$SERVICE_USER}
WorkingDirectory={$REPO_DIR}
Environment=NODE_ENV=production
# Bind to loopback only — Caddy is the sole public ingress on 80/443.
Environment=HOST=127.0.0.1
Environment=PORT={$STOREAI_PORT}
Environment=PATH={$SERVICE_PATH}
ExecStart={$PNPM_BIN} start
# Always restart — the app should survive crashes, OOMs, and DB/Redis blips.
Restart=always
RestartSec=5
# If it keeps crashing fast, back off. StartLimitBurst above caps the rate.
TimeoutStartSec=60
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

# --- Hardening ---------------------------------------------------------------
# Filesystem
ProtectSystem=strict
ProtectHome=true
# The repo lives under /home/<deploy-user>/, which ProtectHome=true would
# normally hide. Make the repo readable, and only .next/cache writable.
BindReadOnlyPaths={$REPO_DIR}
ReadWritePaths={$REPO_DIR}/.next/cache /var/lib/storeai
PrivateTmp=true
PrivateDevices=true

# Kernel + namespace surface
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectProc=invisible
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Network / capabilities / syscalls
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallArchitectures=native

# NOTE: MemoryDenyWriteExecute=true is intentionally commented out.
# V8 (Node.js's JIT) requires RWX memory pages to compile JavaScript; with
# MDWX enabled the process crashes on first JIT compile. If your workload
# runs successfully with --jitless or a future Node removes the need for
# RWX, uncomment the line below.
# MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
